# SOIITP3/resources/tp3_servers.service
[Unit]
Description= tp3 servers

[Service]
Type=simple
ExecStart=/usr/bin/tp3_servers

[Install]
WantedBy=multi-user.target
