# SOIITP3/resources/tp3_users.service
[Unit]
Description= tp3 users

[Service]
Type=simple
ExecStart=/usr/bin/tp3_user
User=tp3user

[Install]
WantedBy=multi-user.target

